# com.johnshelest.robinhoodtrader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.johnshelest.robinhoodtrader</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/johnshelest/Code/robinhood-trader/run.sh</string>
  </array>

  <!-- Run every 5 minutes -->
  <key>StartInterval</key>
  <integer>300</integer>

  <!-- Keep alive on crash -->
  <key>KeepAlive</key>
  <false/>

  <!-- Logging -->
  <key>StandardOutPath</key>
  <string>/Users/johnshelest/Code/robinhood-trader/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnshelest/Code/robinhood-trader/launchd.log</string>

  <!-- Run even on battery -->
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
